# share/systemd/nsd-dnssec-process-resume.timer
[Unit]
Description=resume pending DNSSEC processes

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
RandomizedDelaySec=1m

[Install]
WantedBy=timers.target
